# threadtest.cron
17 12 15 3 * /home/Stefanie/NestedSampling/BeamRecoil/threadtest_ocl_gpu_new.sh
