# threadtest.cron
00 15 16 3 * /home/Stefanie/NestedSampling/BeamRecoil/threadtest_ocl_gpu_new.sh
